# bootstrap-aarch64/var/service/telnetd/run
#!/data/data/com.termux/files/usr/bin/sh
exec busybox telnetd -F
